# deploy/systemd/vps-agent.service
[Unit]
Description=VPS Agent (Slack + OpenCode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-agent
EnvironmentFile=/opt/vps-agent/config/secrets.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.bun/bin
ExecStartPre=/usr/bin/test -x /root/.bun/bin/bun
ExecStart=/root/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
